# packaging/systemd/asgard.service
[Unit]
Description=Asgard control plane for AI/agent development
Documentation=https://github.com/asgard/asgard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/asgard serve
Restart=on-failure
RestartSec=2
# SQLite default — zero external services. Switch to a postgres:// URL to scale out.
Environment=ASGARD_DATABASE_URL=sqlite:////var/lib/asgard/asgard.db
Environment=ASGARD_BIND=0.0.0.0:8080
# Provide a Git token (and any provider keys) via an override or EnvironmentFile.
# EnvironmentFile=/etc/asgard/asgard.env
DynamicUser=yes
StateDirectory=asgard
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes

[Install]
WantedBy=multi-user.target
