# packaging/systemd/frontkeep.service
[Unit]
Description=Frontkeep control plane for AI/agent development
Documentation=https://github.com/glemmestad/frontkeep
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/frontkeep serve
Restart=on-failure
RestartSec=2
# SQLite default — zero external services. Switch to a postgres:// URL to scale out.
Environment=FRONTKEEP_DATABASE_URL=sqlite:////var/lib/glemmestad/frontkeep.db
Environment=FRONTKEEP_BIND=0.0.0.0:8080
# Provide a Git token (and any provider keys) via an override or EnvironmentFile.
# EnvironmentFile=/etc/glemmestad/frontkeep.env
DynamicUser=yes
StateDirectory=frontkeep
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes

[Install]
WantedBy=multi-user.target
